# configs/systemd/cliproxy-api.service.d/10-selfupdate.conf
[Service]
ExecStartPre=-__SELFUPDATE_BIN__
