# configs/systemd/cliproxy-api.service.d/20-rollback.conf
[Unit]
OnFailure=cliproxy-api-rollback.service
StartLimitIntervalSec=300
StartLimitBurst=3
