# configs/systemd/proxybridge.service
[Unit]
Description=ProxyBridge
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash -c '. /usr/local/lib/proxybridge/gen-args.sh && eval exec /usr/local/bin/ProxyBridge $PROXYBRIDGE_ARGS'
ExecStop=/usr/local/bin/ProxyBridge --cleanup
Restart=on-failure
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
